# venditt-player.service
[Unit]
Description=Venditt Player (mpv signage)
After=network-online.target ssh.service
Wants=network-online.target ssh.service

[Service]
Type=simple
User=vobox
Group=vobox
SupplementaryGroups=video render audio input
ExecStartPre=/bin/sleep 15
TimeoutStartSec=30

WorkingDirectory=/data/player

# Start your real script
ExecStart=/bin/bash -lc 'cd /data/player && export GIT_TERMINAL_PROMPT=0 && timeout 20s git pull --ff-only || true; exec /usr/bin/bash ./tvads.sh'

# Auto-restart on crashes
Restart=always
RestartSec=2

KillMode=control-group
TimeoutStopSec=20

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
